# deploy/talos-ui.service
[Unit]
Description=Talos web UI (read-only dashboard)
After=network-online.target talos.service
Wants=network-online.target
# talos.service is a soft dependency — UI can run independently if you want

[Service]
Type=simple
User=talos
Group=talos
EnvironmentFile=-/etc/talos/talos.env
ExecStart=/usr/local/bin/talos-ui --config /etc/talos/talos.yaml
Restart=on-failure
RestartSec=5s

# Hardening (mirrors talos.service)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/talos
StateDirectory=talos
StateDirectoryMode=0750
LogsDirectory=talos
LogsDirectoryMode=0750

# Resource limits
MemoryMax=128M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
