# deploy/talos.service
[Unit]
Description=Talos personal RF security monitor
After=network-online.target kismet.service
Wants=network-online.target
# Kismet is optional — talos can run with a fixture for testing

[Service]
Type=simple
User=talos
Group=talos
EnvironmentFile=-/etc/talos/talos.env
ExecStart=/usr/local/bin/talos --config /etc/talos/talos.yaml
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/talos
StateDirectory=talos
StateDirectoryMode=0750
LogsDirectory=talos
LogsDirectoryMode=0750

# Resource limits (Pi-friendly)
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
